# etc/supervisor/conf.d/x11vnc.conf
[program:x11vnc]
command=bash -l -c "/usr/bin/run-x11vnc.sh"
environment=HOME=/home/user,USER=user
directory=/home/user
autorestart=true
stdout_logfile=/var/log/supervisor/x11vnc.log
redirect_stderr=true
